# deployement/supervisor-example.conf
[program:ics2web]
command = {0}/gunicorn_start
user = {1}
stdout_logfile = {2}
redirect_stderr = true
